# build/fdroidserver.conf
[program:fdroidserver]
user = nobody
command = bash -c 'cd /config && source /usr/local/lib/venv/bin/activate && exec /usr/local/lib/venv/bin/fdroid ${FDROID_COMMAND}'
umask = 000
autorestart = false
startsecs = 0
startretries = 1

[eventlistener:processes]
command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL