# build/fdroidserver.conf
[program:fdroidserver]
user = nobody
command = /home/nobody/start.sh
umask = 000
autorestart = false
startsecs = 0
startretries = 1
stopasgroup = true
killasgroup = true
stopsignal = INT